# var/lib/dpkg/info/ofono.conffiles
/etc/dbus-1/system.d/dundee.conf
/etc/dbus-1/system.d/ofono.conf
/etc/default/ofono
/etc/init.d/dundee
/etc/init.d/ofono
/etc/ofono/phonesim.conf
